# configs/supervisor/kafka_mirror_maker.conf
[program:kafka]
command=/usr/bin/start-kafka_mirror_maker.sh
autostart=true
autorestart=true
